# ansible/roles/magicmirror/files/magicmirror-wake.service
[Unit]
Description=Wake display for MagicMirror morning window
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/johnb/.local/bin/magicmirror-display-wake.sh
