# ansible/roles/magicmirror/files/magicmirror-sleep.timer
[Unit]
Description=Sleep display at 9:00 AM after MagicMirror morning window

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
