# ansible/roles/magicmirror/files/magicmirror-kiosk-wake.service
[Unit]
Description=Wake MagicMirror kiosk — start X server and switch VT
After=graphical.target

[Service]
Type=oneshot
ExecStart=/home/johnb/.local/bin/magicmirror-kiosk-wake.sh
StandardOutput=journal
StandardError=journal
