# ansible/roles/magicmirror/files/magicmirror-kiosk-sleep.service
[Unit]
Description=Sleep MagicMirror kiosk — stop X server and power off monitor
After=graphical.target

[Service]
Type=oneshot
ExecStart=/home/johnb/.local/bin/magicmirror-kiosk-sleep.sh
StandardOutput=journal
StandardError=journal
