# ansible/roles/magicmirror/files/magicmirror-kiosk-sleep.timer
[Unit]
Description=Sleep MagicMirror kiosk at 08:00 daily

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
