# ansible/roles/magicmirror/files/magicmirror-sleep.service
[Unit]
Description=Sleep display after MagicMirror morning window
After=graphical-session.target

[Service]
Type=oneshot
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/johnb/.Xauthority
ExecStart=/bin/bash -c 'xset +dpms && xset dpms 30 30 30 && xset dpms force off'
