# ansible/roles/magicmirror/files/magicmirror-snapshot-cache-system.service
[Unit]
Description=Cache camera snapshots from go2rtc
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=johnb
Group=johnb
ExecStart=/home/johnb/.local/bin/magicmirror-snapshot-cache.sh
StandardOutput=journal
StandardError=journal
