# ansible/roles/magicmirror/files/magicmirror-snapshot-cache.service
[Unit]
Description=Cache camera snapshots from go2rtc

[Service]
Type=oneshot
ExecStart=/home/johnb/.local/bin/magicmirror-snapshot-cache.sh
StandardOutput=journal
StandardError=journal
