# ansible/roles/magicmirror/files/magicmirror-kiosk.service
[Unit]
Description=MagicMirror kiosk — standalone X server with Chromium on VT8
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=johnb
Group=johnb
PAMName=login
ExecStart=/home/johnb/.local/bin/magicmirror-kiosk-x11.sh
Restart=on-failure
RestartSec=5

# X server needs TTY access
TTYPath=/dev/tty8
StandardInput=tty
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
